# templates/sedir.service.j2
[Unit]
Description=Social Enterprise Directory server

[Service]
ExecStart=/usr/bin/node app.js
WorkingDirectory=/sedir/app
Restart=always
RestartSec=20
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=sedir
User=sedir
Group=sedir
Environment=NODE_ENV=production IP_ADDRESS={{ ansible_default_ipv4.address }}

[Install]
WantedBy=multi-user.target
